# build/install/rpm/systemd/modules/appserver-kafka.service
[Unit]
Description=kafka

[Service]
PIDFile=/run/appserver/appserver-kafka.pid
Type=simple
user=root
WorkingDirectory=/root/kafka_2.12-2.5.0/
ExecStart=/root/kafka_2.12-2.5.0/bin/kafka-server-start.sh /root/kafka_2.12-2.5.0/config/server.properties
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target