# build/install/rpm/systemd/modules/appserver-thumbnails.service
[Unit]
Description=AppServer-thumbnails

[Service]
PIDFile=/run/appserver/appserver-thumbnails.pid
Type=simple
user=root
WorkingDirectory=/services/thumb/service/
ExecStart=/usr/bin/dotnet /services/thumb/service/ASC.Thumbnails.Svc.dll --urls=http://0.0.0.0:5016 \
--pathToConf=/app/onlyoffice/config/ --$STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=thumbnails
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target
