# build/install/rpm/systemd/modules/appserver-api.service
[Unit]
Description=AppServer-api

[Service]
PIDFile=/run/appserver/appserver-api.pid
Type=simple
user=root
WorkingDirectory=/var/www/studio/api/
ExecStart=/usr/bin/dotnet /var/www/studio/api/ASC.Web.Api.dll --urls=http://0.0.0.0:5000 \
--pathToConf=/app/onlyoffice/config/ --$STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=api --ENVIRONMENT=test
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target