# build/install/rpm/systemd/modules/appserver-studio.service
[Unit]
Description=AppServer-studio

[Service]
PIDFile=/run/appserver/appserver-studio.pid
Type=simple
user=root
WorkingDirectory=/var/www/studio/server/
ExecStart=/usr/bin/dotnet /var/www/studio/server/ASC.Web.Studio.dll --urls=http://0.0.0.0:5003 \
--pathToConf=/app/onlyoffice/config/ --$STORAGE_ROOT=/app/onlyoffice/data/ --log:dir=/var/log/onlyoffice --log:name=web --ENVIRONMENT=test
TimeoutSec=300
Restart=always

[Install]
WantedBy=multi-user.target